# debian/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable ping-monitor.service
systemctl start ping-monitor.service
